# web/rootfs/etc/services.d/backend/run
#!/usr/bin/with-contenv bash
echo "Starting Backend"

cd /lodestone/backend
rm -rf node_modules
npm install
nodemon
