# storage/rootfs/etc/services.d/minio/run
#!/usr/bin/with-contenv bash

/usr/bin/docker-entrypoint.sh gateway nas /data
