# processor/rootfs/etc/services.d/document-processor/run
#!/usr/bin/with-contenv bash

ls -alt /usr/bin/lodestone-document-processor
/usr/bin/lodestone-document-processor start --redis-hostname=queue
