# elasticsearch/rootfs/etc/services.d/elasticsearch/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid elasticsearch /usr/local/bin/docker-entrypoint.sh
