# processor/rootfs/etc/services.d/processor/run
#!/usr/bin/with-contenv bash

ls -alt /usr/bin/lodestone-processor
/usr/bin/lodestone-processor start --redis-hostname=queue
